# macos/net.geraldhofbauer.spacesaver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Sample LaunchAgent: launch SpaceSaver.app at login so it sits in the menu bar
  as the idle-watcher. Pick the idle delay from its menu-bar icon (1–30 min) —
  the choice persists, no need to edit this file.

  Install:
    1. Copy SpaceSaver.app to /Applications (or ~/Applications and fix the path).
    2. cp macos/net.geraldhofbauer.spacesaver.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/net.geraldhofbauer.spacesaver.plist
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/net.geraldhofbauer.spacesaver.plist
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>net.geraldhofbauer.spacesaver</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/SpaceSaver.app/Contents/MacOS/SpaceSaver</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<!-- No ProcessType on purpose: Background/Adaptive throttle WKWebView's child
	     render processes (measured ~4 fps vs ~20 fps), so leave it at the default
	     Standard QoS band. Idle-watching is a cheap 0.5 s timer regardless. -->
</dict>
</plist>
