# macos/net.geraldhofbauer.spacesaver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Sample LaunchAgent: run SpaceSaver.app as an idle-watcher at login, so it
  behaves like a real screensaver (shows after N seconds idle, hides on input).

  Install:
    1. Copy SpaceSaver.app to /Applications (or ~/Applications and fix the path).
    2. cp macos/net.geraldhofbauer.spacesaver.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/net.geraldhofbauer.spacesaver.plist
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/net.geraldhofbauer.spacesaver.plist

  Adjust the "300" below to the idle delay in seconds you want.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>net.geraldhofbauer.spacesaver</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/SpaceSaver.app/Contents/MacOS/SpaceSaver</string>
		<string>--watch</string>
		<string>300</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
